# config/supervisord/minio.conf
[program:minio]
command=/usr/bin/minio server /mnt/tmpfs/minio
numprocs=1
autostart=true
autorestart=true
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
